# scripts/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable energieip-sol200-authentication.service
systemctl restart energieip-sol200-authentication.service
